# roles/uberspace_mail_pw_reset/templates/flask.ini.j2
[program:uberspace_mail_pw_reset_flask]
directory={{ project_path }}
command={{ project_path }}/bin/uwsgi {{ project_path }}/uwsgi.ini
